# less-code/version-1/crontab
@reboot /usr/bin/python /path/to/nfs-keep.py
